# deploy/supervisor.desktop.ini
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html

; 这个文件只是个配置模板，使用时请从此拷贝出supervisor.ini.
; <virtualenv_path>占位是Python虚拟环境目录, 用实际替换之;
; stdout_logfile 是默认日志目录, 用实际替换之;

[unix_http_server]
file = /tmp/supervisor.sock

[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = /tmp/supervisor.pid
nodaemon = false
minfds = 1024
minprocs = 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock


[program:dockers]
command=<virtualenv_path>/bin/python discuzx_tools/tools/affix_dockers.py
directory=<project_path>/discuzx-tools/
autostart=false
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
numprocs_start=01
